# packaging/rpm/lhsmd.conf
description "Lustre HSM Agent"

start on runlevel [23]
stop on runlevel [S016]

respawn

# Hmm. Do we need a --logfile arg?
exec /usr/sbin/lhsmd > /var/log/lhsmd.log 2>&1
